# dockerfiles/php/crontab
* * * * * php /var/www/html/artisan schedule:run >> /var/log/cron/cron.log 2>&1
# Empty line